# asus-screen-toggle/usr/lib/systemd/user/asus-screen-toggle.service
[Unit]
Description=Asus Screen Toggle User Agent
# Zajistí, že se služba spustí až ve chvíli, kdy běží grafika
After=graphical-session.target
PartOf=graphical-session.target
Documentation=man:asus-user-agent(1)

[Service]
Type=simple
# KLÍČOVÁ OPRAVA: Načte proměnné sezení (WAYLAND_DISPLAY, XDG_RUNTIME_DIR atd.)
# Tyto proměnné tam KDE/Systemd pošle při startu sezení
ImportEnvironment=DISPLAY XAUTHORITY WAYLAND_DISPLAY XDG_RUNTIME_DIR

# Doporučuji přidat parametr -u pro Python, aby logy v journalu neběhaly se zpožděním
ExecStart=/usr/bin/python3 -u /usr/bin/asus-user-agent

KillSignal=SIGTERM
TimeoutStopSec=5
Restart=on-failure
RestartSec=1

# Prostředí GUI
Environment=GTK_USE_PORTAL=1
# Pokud používáte Wayland v KDE, toto může pomoci stabilitě aplikací v Pythonu/GTK
Environment=GDK_BACKEND=wayland,x11

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
