# asus-screen-toggle/DEBIAN/postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

CONFIG_FILE="/etc/asus-check-keyboard.cfg"

# helper: pokud prázdné tak fallback
fallback_if_empty() {
    val="$1"
    fallback="$2"
    if [ -z "$val" ] || [ "$val" = "none" ]; then
        echo "$fallback"
    else
        echo "$val"
    fi
}

case "$1" in
configure)
    # --- pouze pokud soubor neexistuje ---
    if [ ! -e "$CONFIG_FILE" ]; then
        # získání hodnot z debconf
        db_get asus-screen-toggle/usb_keyboard
        SELECTED="$RET"
        VENDOR_ID="0b05"
        PRODUCT_ID="1bf2"
        VENDOR_ID=$(echo "$SELECTED" | cut -d: -f1)
        PRODUCT_ID=$(echo "$SELECTED" | cut -d: -f2 | cut -d' ' -f1)

        db_get asus-screen-toggle/primary_display
        PRIMARY_DISPLAY_NAME=$(fallback_if_empty "$RET" "eDP-1")
        db_get asus-screen-toggle/secondary_display
        SECONDARY_DISPLAY_NAME=$(fallback_if_empty "$RET" "eDP-2")
        db_get asus-screen-toggle/lid
        LID=$(fallback_if_empty "$RET" "LID")

        if [ "$PRIMARY_DISPLAY_NAME" = "$SECONDARY_DISPLAY_NAME" ]; then
          echo "Upozornění: Primární a sekundární výstup jsou stejné ($PRIMARY_DISPLAY_NAME)" >&2
        fi

        # vytvoření konfigu jen při první instalaci/není-li soubor
        install -d -m 0755 "$(dirname "$CONFIG_FILE")"
        cat <<EOF > "$CONFIG_FILE"
VENDOR_ID="$VENDOR_ID"
PRODUCT_ID="$PRODUCT_ID"
PRIMARY_DISPLAY_NAME="$PRIMARY_DISPLAY_NAME"
SECONDARY_DISPLAY_NAME="$SECONDARY_DISPLAY_NAME"
LID="$LID"
EOF
        # zaregistruj conffile (volitelné, užitečné když soubor vzniká v postinst)
        if command -v dpkg-maintscript-helper >/dev/null 2>&1; then
            dpkg-maintscript-helper addconffile "$CONFIG_FILE"
        fi
    else
        # soubor již existuje – zachováme ho, pokud uživatel upravoval
        echo "Konfig $CONFIG_FILE již existuje, nebude přepsán."
    fi

    /usr/bin/asus-check-keyboard-genrules.sh

    if command -v udevadm >/dev/null; then
      udevadm control --reload-rules
      udevadm trigger
    fi

    if [ "$(ps -p 1 -o comm=)" = "systemd" ]; then
      /usr/bin/systemctl enable asus-bottom-screen-init.service
      /usr/bin/systemctl start asus-bottom-screen-init.service
    else
       echo "systemd not running, enable service later."
    fi
    ;;
*)
    # other states: upgrade, abort, etc. Nic dělat
    ;;
esac

#DEBHELPER#
exit 0
