# asus-screen-toggle/usr/share/asus-screen-toggle/autostart/asus-screen-toggle.desktop
[Desktop Entry]
Type=Application
Name=asus-screen-toggle
Exec=/usr/local/bin/asus-check-keyboard-user.sh
X-GNOME-Autostart-enabled=true
